# general/supervisor/conf.d/pf-trilium.conf
[program:pf-trilium]
command=/Users/chloelee/bin/pf-trilium
autostart=true
autorestart=true
startretries=5
